# images/archivesspace/rootfs/etc/s6-overlay/s6-rc.d/archivesspace/run
#!/usr/bin/env bash
set -e
exec /usr/local/bin/archivesspace-startup.sh
