# images/islandora/rootfs/etc/s6-overlay/s6-rc.d/set-site-defaults/up
/etc/s6-overlay/scripts/set-site-defaults.sh
