# images/base/rootfs/etc/s6-overlay/s6-rc.d/vault-secrets/up
/etc/s6-overlay/scripts/vault-secrets.sh
